# .config/systemd/user/bedtime.service
[Unit]
Description=Bedtime notification service

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "Go to bed"
