# .config/systemd/user/lockout.service
[Unit]
Description=After-hours i3 lockout service

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /tmp/i3-lockout
