# .config/systemd/user/taskwarrior-notify.service
[Unit]
Description=TaskWarrior synchronization service

[Service]
Type=oneshot
ExecStart=%h/.bin/taskwarrior-notify
